# tool-labs/crontab
PATH=/usr/local/bin:/usr/bin:/bin

# Crontab entries for SuggestBot

# Redirects stderr to an error file in our home dir to catch errors,
# and redirect stdout to /dev/null because the shell scripts
# are set up to redirect stderr and stdout to files themselves.

# Jobs updating inlink count tables
# Run every morning for enwiki, nowiki, svwiki, ptwiki, fawiki, huwiki,
# and ruwiki
14 1 * * * qsub -q task $HOME/link-rec/update_linkcounts.sh no 2>$HOME/crontab-errors.txt >/dev/null
14 2 * * * qsub -q task $HOME/link-rec/update_linkcounts.sh sv 2>$HOME/crontab-errors.txt >/dev/null
14 3 * * * qsub -q task $HOME/link-rec/update_linkcounts.sh pt 2>$HOME/crontab-errors.txt >/dev/null
14 4 * * * qsub -q task $HOME/link-rec/update_linkcounts.sh fa 2>$HOME/crontab-errors.txt >/dev/null
14 5 * * * qsub -q task $HOME/link-rec/update_linkcounts.sh en 2>$HOME/crontab-errors.txt >/dev/null
14 6 * * * qsub -q task $HOME/link-rec/update_linkcounts.sh ru 2>$HOME/crontab-errors.txt >/dev/null
14 7 * * * qsub -q task $HOME/link-rec/update_linkcounts.sh fr 2>$HOME/crontab-errors.txt >/dev/null

# Job(s) updating en:Wikipedia:Community portal/Opentask
15,45 * * * * qsub -q task $HOME/projects/opentask/opentask.sh 2>$HOME/crontab-errors.txt >/dev/null

# Job(s) updating lists of open tasks on other Wikis
31 * * * * qsub -q task $HOME/projects/opentask/launch-lang.sh pl "Wikipedia:Sprzątanie Wikipedii/Sugerowane listy" $HOME/projects/opentask/plwiki-taskdef.json 2>$HOME/crontab-errors.txt >/dev/null
25 * * * * qsub -q task $HOME/projects/opentask/launch-lang.sh ru "Википедия:Сообщество/Задачи" $HOME/projects/opentask/ruwiki-taskdef.json 2>$HOME/crontab-errors.txt >/dev/null
5 2 2-30/2 * * qsub -q task $HOME/projects/opentask/launch-lang.sh th "วิกิพีเดีย:ศาลาประชาคม/Opentask" $HOME/projects/opentask/thwiki-taskdef.json

# Job(s) publishing inlink counts for Wikis
# 35 7 * * 7 qsub -q task $HOME/projects/inlink-counts/publish-inlink-counts.sh en 2>$HOME/crontab-errors.txt >/dev/null
